# ddns-updater.service.tmpl
#encoding utf-8

#raw
################################################################################
# DO NOT MODIFY THIS FILE.                                                     #
# It is generated by the template named like this file and a .tmpl suffix      #
# added. Make modifications to the template and regenerate this file.          #
################################################################################
#end raw

[Unit]
Description=Reflect changes in the IP assigned by ISP in the strato DDNS service

[Service]
Type=simple
ExecStart=/usr/local/bin/$app_name $mode -loop -debug
# `-daemon` might be necessary in order to avoid `Service hold-off time over, scheduling restart.` of `systemd`, but using `-daemon` and `Type=forking` causes `systemd` unit to fail for no reason
\# the service is very likely to fail, e.g. when a timeout or "temporary failure in name resolution" occurs
Restart=on-failure
RestartSec=15s
User=$user
Group=$group

[Install]
WantedBy=multi-user.target
# might avoid `Service hold-off time over, scheduling restart.`<ref>http://unix.stackexchange.com/questions/213185/restarting-systemd-service-on-dependency-failure</ref>
